# schedule.cron
*/30 * * * * cd /Volumes/SeagateD/Projects/SearchTrendsSocialMediaData; . venv/bin/activate; ./trends.py